# contentpipe.service
[Unit]
Description=ContentPipe — AI 图文内容生产线
After=network.target clawdbot.service
Wants=clawdbot.service

[Service]
Type=simple
User=ssp
WorkingDirectory=/home/ssp/work/openclawWS/clawdbot_workspace/plugins/content-pipeline/scripts
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3 -m uvicorn web.app:app --host 0.0.0.0 --port 8765
Environment="PYTHONPATH=."
Environment="PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin:/home/ssp/.local/bin:/home/ssp/.npm-global/bin"
EnvironmentFile=-/home/ssp/work/openclawWS/clawdbot_workspace/plugins/content-pipeline/.env.local
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
